# .runit/.runit/sv/psd/run
#!/bin/sh

psd startup
psd resync
exec chpst -b psd pause
