# .runit/.runit/sv/uxplay/run
#!/bin/sh
exec 2>&1
exec uxplay -n $(hostname) -nh -bt709 -nohold -vs ximagesink -fps 60 -s 1920x1080@60
