# .runit/.runit/sv/uxplay/run
#!/bin/sh
exec 2>&1
exec uxplay -p -fps 60 -s 1920x1080@60
