# .runit/sv/dunst/run
#!/bin/sh

exec dunst -conf /home/miuka/.config/dunst/dunstrc 2>&1
