# .runit/sv/picom/run
#!/bin/sh

exec picom --dbus --experimental-backends 2>&1
