# .runit/.runit/sv/rustdesk/run
#!/bin/sh

exec /usr/bin/rustdesk --service 2>&1
