# .runit/sv/smartdns/run
#!/bin/sh

exec smartdns -p /run/user/$(id -u)/smartdns -c ~/.config/smartdns/smartdns.conf -f 2>&1
