# .runit/sv/smartdns/run
#!/bin/sh

exec smartdns -p - -c ~/.config/smartdns/smartdns.conf -f 2>&1
