# .runit/.runit/sv/spotifyd/run
#!/bin/sh

exec spotifyd --no-daemon --config-path $HOME/.config/spotifyd/spotifyd.conf >&1
