# .runit/.runit/sv/mpd-notify/run
#!/bin/sh

# Wait for the mpd.
pgrep "mpd" > /dev/null || exit 1
exec mpd-notify 2>&1
