# .runit/.runit/sv/avahi-daemon/run
#!/bin/sh
# mv it to /etc/runit/sv/
exec 2>&1
sv check dbus >/dev/null || exit 1
exec avahi-daemon
