# .runit/sv/easyeffects/run
#!/bin/sh

exec easyeffects --gapplication-service 2>&1
