# .runit/sv/easyeffects/run
#!/bin/sh

easyeffects --gapplication-service 2>&1
