# .runit/.runit/sv/mega/run
#!/bin/sh

exec mega-cmd-server --skip-lock-check 2>&1
