# .runit/.runit/sv/frpc/run
#!/bin/sh

exec /usr/bin/frpc -c ~/.config/frpc/frpc.toml 2>&1
