# .runit/.runit/sv/fusuma/run
#!/bin/sh

exec fusuma -c $HOME/.config/fusuma/config.yml 2>&1
