# .runit/sv/personalDNSfilter/run
#!/bin/sh

exec /usr/bin/java -cp /usr/local/lib/personalDNSfilter/personalDNSfilter.jar dnsfilter.DNSFilterProxy 2>&1
